# systemd/solver-lnk.service
# Example systemd service for solver-lnk
# Copy to /etc/systemd/system/solver-lnk.service and adjust paths
#
# Install:
#   sudo systemctl daemon-reload
#   sudo systemctl enable solver-lnk.service
#   sudo systemctl start solver-lnk.service
#
# Status:
#   systemctl status solver-lnk --no-pager
#   journalctl -u solver-lnk -f

[Unit]
Description=Lords and Knights Solver (gRPC)
After=network.target

[Service]
Type=simple

# TODO: Change to your user/group
User=YOUR_USER
Group=YOUR_USER

# TODO: Change to your repo path
WorkingDirectory=/path/to/solver-lnk
ExecStart=/path/to/solver-lnk/server

Restart=always
RestartSec=5

# Localhost-only networking
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Logging
StandardOutput=journal
StandardError=journal

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
